# server/systemd/ink2ai_cleanup.timer
[Unit]
Description=Timer to delete files older than 24 hours every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
